# files/system/usr/lib/systemd/user/distrobox-auto-pull.timer
[Unit]
Description=Distrobox auto-update timer

[Timer]
OnCalendar=hourly
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
